# openRC-initscripts/joycond
#!/sbin/openrc-run
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

depend(){
	need udev
	need bluetoothd
}
start(){
	ebegin "Joycond Userspace Driver"
	start-stop-daemon --start /usr/bin/joycond --make-pidfile \
		--background --pidfile /var/run/joycond.pid
	eend $?
}
stop(){
	ebegin "Stopping Joycond Userspace Driver"
	start-stop-daemon --stop --pidfile /var/run/joycond.pid
}
