# openRC-initscripts/r2p
#!/sbin/openrc-run
# Copyright 2019 Alexander Weber (bell07)
#Distributed under the terms of the GNU General Public License v2

depend(){
	need sysfs
}
start(){
	ebegin "Enable reboot2payload with default_payload"
	# Kernel 4.9.201, kernel 4.8.140
	echo 1 > /sys/bus/platform/drivers/tegra-reboot2payload/r2p/default_payload_ready
	# if not, check /sys/bus/platform/drivers/tegra-reboot2payload/r2p/ folder	
			eend $?
}
