# openRC-initscripts/switchroot-pre-config
#!/sbin/openrc-run
# Copyright 2021 atlasflamer
# Distributed under the terms of the GNU General Public License v2
command=/usr/bin/switchroot-pre-config
depend(){
	need net
}

